# deployment-aws/puppet/manifests/cromwell-supervisor.conf.epp
[program:cromwell]

command = /usr/bin/java -Dconfig.file=/etc/cromwell.conf -jar /opt/cromwell.jar server

user = <%= $cromwell_user %>

directory = /

; Put process stdout output in this file
stdout_logfile = <%= $cromwell_log_dir %>/stdout.log

; Put process stderr output in this file
stderr_logfile = <%= $cromwell_log_dir %>/stderr.log

autorestart=true
stopsignal=QUIT
