# deploy/puppet/mevapi/templates/supervisor/cloud_sql_proxy.conf.epp
[program:cloud_sql_proxy]

command=/opt/software/database/cloud_sql_proxy -instances=%(ENV_DB_HOST_FULL)s -dir=%(ENV_CLOUD_SQL_MOUNT)s

user=<%= $mevapi::app_user %>

; Put process stdout output in this file
stdout_logfile = <%= $mevapi::log_dir %>/%(program_name)s.log

; Put process stderr output in stdout file
redirect_stderr = true

autostart=false
autorestart=true
stopsignal=QUIT
